# templates/sonar.service.erb
[Unit]
Description=Sonar <%= @version.to_s %>
After=network.target network-online.target
Wants=network-online.target

[Service]
LimitNOFILE=65536
LimitNPROC=4096
ExecStart=<%= @installroot.to_s -%>/sonar/bin/<%= @arch.to_s.gsub(/_/,'-') -%>/sonar.sh start
ExecStop=<%= @installroot.to_s -%>/sonar/bin/<%= @arch.to_s.gsub(/_/,'-') -%>/sonar.sh stop
ExecReload=<%= @installroot.to_s -%>/sonar/bin/<%= @arch.to_s.gsub(/_/,'-') -%>/sonar.sh restart
PIDFile=<%= @real_home.to_s -%>/sonar.pid
Type=forking
User=<%= @user.to_s %>

[Install]
WantedBy=multi-user.target
